# AIM/_archive/v7_ai_code/deploy/systemd/aim-generalist.service
[Unit]
Description=AIM generalist agent (Rust :8774)
After=network-online.target aim-llm.service
Wants=aim-llm.service

[Service]
Type=simple
User=oem
WorkingDirectory=/home/oem/Desktop/LC/AIM/rust-core
EnvironmentFile=-/home/oem/.aim_env
Environment="RUST_LOG=info"
Environment="AIM_ENV=prod"
Environment="AIM_GENERALIST_ROOT=/home/oem/Desktop/LC/AIM/Patients"
ExecStart=/home/oem/Desktop/LC/AIM/rust-core/target/release/aim-generalist
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/oem/Desktop/LC/AIM/Patients
PrivateTmp=true

[Install]
WantedBy=multi-user.target
