# AIM/_archive/v7_ai_code/deploy/systemd/aim-weekly-digest.timer
[Unit]
Description=Run AIM weekly digest every Sunday at 09:00
Requires=aim-weekly-digest.service

[Timer]
OnCalendar=Sun *-*-* 09:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
