# AIM/_archive/v7_ai_code/deploy/systemd/aim-hive-worker.user.service
[Unit]
Description=AIM Hive worker (USER-MODE) — periodic anonymized telemetry contribution
Documentation=file:///home/oem/Desktop/LC/AIM/AI/HIVE_ARCHITECTURE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/oem/Desktop/LC/AIM
EnvironmentFile=-/home/oem/.aim_env
Environment="RUST_LOG=info"
ExecStart=/home/oem/Desktop/LC/AIM/rust-core/target/release/aim-hive-telemetry contribute
ExecStartPost=/usr/bin/touch /home/oem/.cache/aim/hive_last_contribute.txt
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
