# AIM/_archive/v7_ai_code/deploy/systemd/aim-weekly-digest.service
[Unit]
Description=AIM weekly self-improvement digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.aim_env
WorkingDirectory=%h/Desktop/LC/AIM
ExecStart=%h/Desktop/LC/AIM/venv/bin/python -m scripts.weekly_digest
StandardOutput=journal
StandardError=journal
