# AIM/deploy/aim-transcriber.service
[Unit]
Description=AIM Telegram Transcriber — голосовые сообщения → текст
After=network.target

[Service]
Type=simple
User=oem
WorkingDirectory=/home/oem/Desktop/LC/AIM
EnvironmentFile=/home/oem/.aim_env
# --lang ka : грузинский
# --fix     : DeepSeek коррекция
# --mic     : автовключение микрофона врача
ExecStart=/home/oem/Desktop/LC/AIM/venv/bin/python /home/oem/Desktop/LC/AIM/bot_transcriber.py --mic
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
